# pi_config_files/pifinder_splash.service
[Unit]
Description=PiFinderSplash

[Service]
WorkingDirectory=/home/pifinder/PiFinder5/python
ExecStart=/home/pifinder/PiFinder5/python/.venv/bin/python -m PiFinder.splash

[Install]
WantedBy=basic.target
